# systemd/autonyaa.service
[Unit]
Description=Autonyaa

[Service]
Type=oneshot
ExecStart=/usr/bin/autonyaa/autonyaa.py

[Install]
WantedBy=multi-user.target
